# target/linux/mvebu/base-files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
echo 1 > /sys/fs/cgroup/memory.use_hierarchy
mkdir -p /sys/fs/cgroup/systemd
mount -t cgroup -o none,name=systemd systemd /sys/fs/cgroup/systemd
exit 0
